# systemd/fpvcopilot-sky.service
[Unit]
Description=FPV Copilot Sky - Backend Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=fpvcopilotsky
Group=fpvcopilotsky
# Ensure we have access to serial ports and video devices
SupplementaryGroups=dialout video
WorkingDirectory=/opt/FPVCopilotSky
Environment="PATH=/opt/FPVCopilotSky/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Use virtual environment Python
ExecStart=/opt/FPVCopilotSky/venv/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fpvcopilot-sky

# Security and resource limits
LimitNOFILE=65536
LimitNPROC=4096

# --- Security hardening ---
# Isolate /tmp so other services cannot access our temp files
PrivateTmp=true
# Prevent access to /home, /root, /run/user
ProtectHome=true
# Make /usr, /boot, /efi read-only (not /opt, which we need writable for preferences)
ProtectSystem=full
# Prevent loading kernel modules
ProtectKernelModules=true
# Prevent modifying control groups
ProtectControlGroups=true
# Deny creating SUID/SGID files
RestrictSUIDSGID=true

# --- Resource limits ---
# Cap memory at 512 MB to prevent OOM on embedded boards
MemoryMax=512M
# Allow full CPU — GStreamer x264enc uses native threads across all 4 cores
# Previous CPUQuota=80% throttled encoding performance
# CPUQuota removed intentionally

# --- CPU governor: lock to 'performance' while service is running ---
# '+' prefix runs as root (needed for sysfs writes) even though service runs as fpvcopilotsky
ExecStartPre=+/bin/sh -c 'for g in /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor; do echo performance > "$$g" 2>/dev/null || true; done'
ExecStopPost=+/bin/sh -c 'for g in /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor; do echo ondemand > "$$g" 2>/dev/null || true; done'

[Install]
WantedBy=multi-user.target
